# root/etc/s6-overlay/s6-rc.d/svc-uwsgi/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

echo "[svc-uwsgi] Starting uwsgi server"

uwsgi --ini /etc/uwsgi/uwsgi.ini