# root/etc/s6-overlay/s6-rc.d/init-db-config/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

echo "[init-db-config] Running database migrations"

# database migration
(cd /the_wiki && python3 manage.py migrate --no-input)

